# root/etc/s6/services/chrome/run
#!/usr/bin/execlineb -P

with-contenv
/usr/bin/chromm
